# services/fluentbit-watcher.path
[Unit]
Description=Restart fluentbit on config changes
StartLimitIntervalSec=0

[Path]
PathModified=/etc/fluent-bit.yaml

[Install]
WantedBy=multi-user.target
