# services/mavlink-watcher.path
[Unit]
Description=Restart mavlink on ENV changes
StartLimitIntervalSec=0

[Path]
PathModified=/etc/mavlink.env

[Install]
WantedBy=multi-user.target
