# services/mavlink-restart.service
[Unit]
Description=One-time restart of mavlink

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart mavlink.service
